# deploy/launchd/com.alpaca-trading.health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alpaca-trading.health</string>
    <key>ProgramArguments</key>
    <array>
        <string>__VENV_PATH__/bin/python</string>
        <string>__PROJECT_PATH__/main.py</string>
        <string>health_check</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__PROJECT_PATH__</string>
    <!-- 30分ごとに実行 -->
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>StandardOutPath</key>
    <string>__PROJECT_PATH__/logs/launchd-health.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_PATH__/logs/launchd-health-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__VENV_PATH__/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
