# deploy/launchd/com.alpaca-trading.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alpaca-trading.morning</string>
    <key>ProgramArguments</key>
    <array>
        <string>__VENV_PATH__/bin/python</string>
        <string>__PROJECT_PATH__/main.py</string>
        <string>morning</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__PROJECT_PATH__</string>
    <!-- 14:25 UTC = 9:25 EST / 10:25 EDT. Python側のis_market_open()でDST補正 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>25</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>25</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>__PROJECT_PATH__/logs/launchd-morning.log</string>
    <key>StandardErrorPath</key>
    <string>__PROJECT_PATH__/logs/launchd-morning-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__VENV_PATH__/bin:__HOME_PATH__/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
